# systemd/kaktusas.timer
[Unit]
Description=kaktusas timer unit

[Timer]
OnCalendar=Mon,Tue,Wed,Thu,Fri 8,10,12,14,16,18
